# cron/crontab
# CRONJOBS
PATH=/usr/local/bin:/usr/bin:/bin
15 2 * * * full_backup_procedure > /proc/1/fd/1 2>/proc/1/fd/2
