# cron/crontab
# Supercronic crontab - runs as non-root user
# Environment variables are inherited from the container
CRON_TZ=Europe/Amsterdam

15 2 * * * source ${HOME}/.backup_env && full_backup_procedure
